# docker/runit/netconfig/run
#!/bin/sh
cd /home/netconfig/
chpst -u netconfig:www-data /usr/bin/uwsgi --ini netconfig.ini
